# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=%(here)s/logs/supervisord.log
pidfile=/tmp/supervisord.pid
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:redis]
command=redis-server --port 6379 --dir /tmp --loglevel warning
autostart=true
autorestart=true
stdout_logfile=%(here)s/logs/redis.log
stderr_logfile=%(here)s/logs/redis_error.log
priority=10

[program:celery]
command=%(here)s/venv/bin/celery -A mmportal worker --loglevel=info --pool=solo
directory=%(here)s
autostart=true
autorestart=true
stdout_logfile=%(here)s/logs/celery.log
stderr_logfile=%(here)s/logs/celery_error.log
priority=20
environment=PATH="%(here)s/venv/bin:%(ENV_PATH)s"

[program:django]
command=%(here)s/venv/bin/python manage.py runserver 8001
directory=%(here)s
autostart=true
autorestart=true
stdout_logfile=%(here)s/logs/django.log
stderr_logfile=%(here)s/logs/django_error.log
priority=30
environment=PATH="%(here)s/venv/bin:%(ENV_PATH)s"
